# cpuscale.service
[Unit]
Description=cpu auto scaling

[Service]
Restart=on-failure
RestartSec=30
ExecStart=/usr/bin/cpuscale

[Install]
WantedBy=multi-user.target
