# cpuscale.service
[Unit]
Description=cpu auto scaling

[Service]
ExecStart=/usr/bin/cpuscale

[Install]
WantedBy=multi-user.target
